# systemd/cert-checker.socket
[Socket]
ListenStream=/tmp/CertCheckerSocket
SocketUser=system
SocketUser=system
SocketMode=0777

Service=cert-checker.service
